# cyrus3/service/cyrus/run
#!/bin/bash
set -e

/usr/lib/cyrus/bin/master -D
