# cyrus3/service/syslog-ng/run
#!/bin/bash
set -e

/usr/sbin/syslog-ng -p /var/run/syslog-ng.pid -F
